# raspberry_pi/dronecan-esc.service
[Unit]
Description=DroneCAN ESC Controller
After=network.target pigpiod.service
Requires=pigpiod.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/home/nicolay
ExecStart=/usr/bin/python3 /home/nicolay/dronecan_esc_controller.py --pwm --quiet
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Sicherheits-Einstellungen
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
