# raspberry_pi/motor_controller_v2.service
[Unit]
Description=Quassel UGV Motor Controller v2.0
After=network.target pigpiod.service can0.service
Requires=pigpiod.service
Wants=can0.service

[Service]
Type=simple
User=nicolay
Group=nicolay
WorkingDirectory=/home/nicolay/motor_controller
Environment="PYTHONUNBUFFERED=1"

# Haupt-Befehl mit Config-Datei
ExecStart=/usr/bin/python3 -m motor_controller.main --config /home/nicolay/motor_controller/config.yaml

# Alternativ: Legacy CLI-Argumente (f체r R체ckw채rtskompatibilit채t)
# ExecStart=/usr/bin/python3 -m motor_controller.main --pwm --pins 18,19 --web --web-port 80 --can can0

# Restart-Verhalten
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Sicherheit
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=motor-controller

[Install]
WantedBy=multi-user.target

